# ideapad-S300/rc.local
#!/bin/sh

## rc.local for ideapad S300
## work with systemd service rc-local
## path /etc/rc.d/rc.local
#
## INSTALL
## chmod 755 /etc/rc.d/rc.local
## chown root:root /etc/rc.d/rc.local
## systemctl enable rc-local

sleep 5
/usr/sbin/ethtool -s enp1s0 autoneg off speed 100
/usr/sbin/ethtool -s enp1s0 wol d

exit 0
