# ideapad-S300/rc.local
#!/bin/sh

# chmod 755 /etc/rc.d/rc.local
# systemctl enable rc-local

sleep 5
powertop --auto-tune 2>/dev/null

exit 0
